# systemd/status.timer
[Unit]
Description=Run status check every 1 minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
Unit=status.service

[Install]
WantedBy=timers.target