# systemd/status.service
[Unit]
Description=Health Check

[Service]
Type=oneshot
ExecStart= ../scripts/status.sh
